# tools/maccontrol/app.openclaw.maccontrol.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>app.openclaw.maccontrol</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/ngotrung/Applications/MacControl.app/Contents/MacOS/MacControl</string>
    <string>serve</string>
    <string>--port</string>
    <string>17891</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/ngotrung/OpenClawTools/maccontrol.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ngotrung/OpenClawTools/maccontrol.err</string>
</dict>
</plist>
